# external/packages/overlays/motd/etc/update-motd.d/10-orangerigol-header
#!/bin/bash

. /etc/os-release
#BOARDNAME=$(cat /etc/hostname)
KERNELID=$(uname -r)

#toilet -f standard  -F metal ${BOARDNAME}
echo -e "\033[0;33m" ; toilet -f standard Orange Rigol ; echo -e "\033[0m"
echo ""
echo -e "Welcome to \e[0;91m${NAME} \x1B[0m$(if [[ $ID=ubuntu ]]; then echo ${UBUNTU_CODENAME^}; fi) with \e[0;91mLinux $KERNELID\x1B[0m\n"
